# docker/demo/s6-rc.d/redis/run
#!/command/with-contenv sh
set -eu

exec s6-setuidgid redis redis-server /etc/shortener-demo/redis.conf
